# meta-openpli/recipes-support/partitions-by-name/files/device-mapper.rules
KERNEL=="mtd*", ENV{DEVTYPE}=="mtd", SYMLINK+="block/by-name/$sysfs{name}"

KERNEL=="mmcblk0p*", SUBSYSTEM=="block", RUN+="/etc/udev/scripts/mmc-dev-mapping.sh"

